# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:monitor]
command=penne_monitor /home/ftpuser

[program:celery_tasks]
command=celery -A penne_shell.tasks worker