# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:monitor]
command=monitor

[program:celery_tasks]
command=celery -A penne_shell.tasks worker
