# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:monitor]
command=monitor
